# com.sk.agent.harvest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sk.agent.harvest</string>
    <key>ProgramArguments</key>
    <array>
        <string>REPO_PATH/bin/sk-harvest-cron</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>0</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>8</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>WorkingDirectory</key>
    <string>REPO_PATH</string>
    <key>StandardOutPath</key>
    <string>REPO_PATH/reports/harvest-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>REPO_PATH/reports/harvest-stderr.log</string>
</dict>
</plist>
